# modules/biometrics-archiver/sleepypod-biometrics-pruner.timer
[Unit]
Description=Run sleepypod biometrics archive pruner every 15 min

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Persistent=true
Unit=sleepypod-biometrics-pruner.service

[Install]
WantedBy=timers.target
