# modules/calibrator/sleepypod-calibrator.service
[Unit]
Description=SleepyPod Calibrator (Sensor Baselines and Thresholds)
Documentation=https://github.com/sleepypod/core/tree/main/modules/calibrator
After=network.target sleepypod.service
Wants=sleepypod.service

[Service]
Type=simple
User=dac
UMask=0002
WorkingDirectory=/opt/sleepypod/modules/calibrator
ExecStart=/opt/sleepypod/modules/calibrator/.venv/bin/python main.py
Restart=always
RestartSec=30

# Environment — override via /etc/sleepypod/modules/calibrator.env
EnvironmentFile=-/etc/sleepypod/modules/calibrator.env
Environment="RAW_DATA_DIR=/persistent"
Environment="DATABASE_URL=file:/persistent/sleepypod-data/sleepypod.db"
Environment="BIOMETRICS_DATABASE_URL=file:/persistent/sleepypod-data/biometrics.db"
Environment="CALIBRATION_HOUR=6"

# Resource limits (calibration loads up to 6h of sensor data into numpy
# arrays for scipy processing — 500Hz piezo × 6h easily exceeds a few
# hundred MB). 256M was tripping systemd's kill at ~6s into capacitance
# calibration on Pod 3 (status=9/KILL). 1G is conservative vs total
# Pod RAM (~6GB) and leaves plenty of headroom.
MemoryMax=1G
CPUQuota=25%

# Sandboxing
NoNewPrivileges=true
# ReadWritePaths overrides ReadOnlyPaths for the sub-path (systemd precedence rules)
ReadOnlyPaths=/persistent
ReadWritePaths=/persistent/sleepypod-data
ProtectSystem=strict
ProtectKernelTunables=true
# Isolated /tmp — otherwise ProtectSystem=strict + narrow ReadWritePaths
# hide host /tmp → Python tempfile.mkstemp() errors "No usable temporary
# directory found". See piezo-processor.service for the full rationale.
PrivateTmp=true

[Install]
WantedBy=multi-user.target
