# modules/biometrics-archiver/frank.service.d/sleepypod-biometrics-tmpfs.conf
# Drop-in for /lib/systemd/system/frank.service installed by sleepypod's
# biometrics-archiver module. Pulls in the tmpfs mount + prep service so
# frankenfirmware finds /persistent/biometrics ready (mounted, with
# symlinks) when frank.sh execs.
[Unit]
Wants=sleepypod-tmpfs-prep.service
After=sleepypod-tmpfs-prep.service
RequiresMountsFor=/persistent/biometrics
