# modules/biometrics-archiver/frank.service.d/sleepypod-biometrics-tmpfs.conf
# Drop-in for /lib/systemd/system/frank.service installed by sleepypod's
# biometrics-archiver module. Pulls in the tmpfs mount + prep service so
# frankenfirmware finds /persistent/biometrics ready (mounted, with
# symlinks) when frank.sh execs.
[Unit]
Wants=sleepypod-tmpfs-prep.service
After=sleepypod-tmpfs-prep.service
RequiresMountsFor=/persistent/biometrics

[Service]
# On Pod variants with frank.sh, the patched shim does `cd /persistent/biometrics`
# itself and this is a no-op (shell cd wins). On newer Pod 5 variants without
# the shim, frank.service ExecStarts frankenfirmware directly — without this
# override, ~1 GB/day of RAW frames lands on /persistent eMMC instead of the
# tmpfs at /persistent/biometrics, wearing the flash.
WorkingDirectory=/persistent/biometrics
